# infra/cron/crontab
0 8 * * * /scripts/run-daily.sh >> /var/log/cron.log 2>&1
